# session-service/supervisord.conf
[supervisord]
nodaemon=true

[program:session-service]
command=/bin/bash /app/utils/start-session-server.sh
